# install/supervisor.conf
[supervisord]
nodaemon=true

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /apisrv/install/hogwhat_uwsgi.ini

[program:nginx]
command=/usr/sbin/nginx